# roles/mail/templates/mail-notify.service
[Unit]
Description=Mail Notification Daemon

[Service]
ExecStart=/usr/local/bin/mail-notify-daemon

[Install]
WantedBy=default.target
